# templates/twemproxy_upstart.conf.erb
description "nutcracker/twemproxy"

env DAEMON=/usr/local/sbin/nutcracker
env DAEMON_ARGS="-c /etc/twemproxy/nutcracker.yml -o /var/log/twemproxy/nutcracker.log &2>1"
env USER=twemproxy

start on startup
stop on shutdown

respawn

exec start-stop-daemon --start --make-pidfile --pidfile /var/run/twemproxy.pid --chuid $USER --exec $DAEMON -- $DAEMON_ARGS